# acme/acme_service/bin/acme
#!/bin/bash
# init file for ACME
#
# description: Ultra*Log Automated Configuration Management Environment (ACME) Daemon
# chkconfig: 2345 99 10
#
# source function library... run chkconfig --add acme to install
. /etc/init.d/functions

start() {
  ps -ef > /tmp/$$  
  if grep "acme\.rb" /tmp/$$> /dev/null ; then
    echo "Acme already running.  Use service acme stop first."
    RETVAL=1
  else
    cd /usr/local/acme
    # su -l acme_user #...defaults to root!
    #/usr/bin/ruby acme.rb &> run.log &

    # restart if it dies
    ./acme.sh &>/dev/null &

    echo "$!" > acme-sh.pid

	RETVAL=0
	echo
	echo_success
  fi
  rm /tmp/$$
  return $RETVAL
}

stop() {
  cd /usr/local/acme
  if [ -f acme-sh.pid ] ; then
    local line p pid=
    read line < acme-sh.pid
    for p in $line ; do
      [ -z "${p//[0-9]/}" -a -d /proc/$p ] &&
      pid="$pid $p"
    done
    if [ -n "${pid-:}" ] ; then
      kill -9 $pid
      rm -rf acme-sh.pid
    fi
  fi  
  if [ -f acme.pid ] ; then
    local line p pid=
    read line < acme.pid
    for p in $line ; do
      [ -z "${p//[0-9]/}" -a -d /proc/$p ] &&
      pid="$pid $p"
    done
    if [ -n "${pid-:}" ] ; then
      kill -9 $pid
      rm -rf acme.pid
    fi
  fi  
  
	RETVAL=0
	echo_success
	return $RETVAL
}

reload(){
	stop
	start
}

restart(){
	stop
	start
}

case "$1" in
  start)
	start
	;;
  stop)
	stop
	;;
  restart)
	restart
        ;;
  reload)
	reload
        ;;
  status)
        status snmpd
	RETVAL=$?
        ;;
  *)
	echo $"Usage: $0 {start|stop|status|restart|reload}"
	RETVAL=1
esac

exit $RETVAL
